# etc/systemd/system/systemd-udevd.service.d/ipaddressallow.conf
# allow curl
[Service]
IPAddressAllow=127.0.0.1
# gnudb.org - audiocd db
IPAddressAllow=78.47.118.0
