# etc/udev/rules.d/rtl-sdr.rules
SUBSYSTEM!="bluetooth", SUBSYSTEM!="net", SUBSYSTEMS=="usb", ENV{ID_SOFTWARE_RADIO}="1", RUN+="/srv/http/bash/settings/features-data.sh pushrefresh"
