# etc/systemd/system/shairport-sync.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/taskset -c 3 /usr/bin/shairport-sync
User=root
Group=root
