# etc/systemd/system/camilladsp.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/camilladsp /srv/http/data/camilladsp/configs/camilladsp.yml --port 1234 --logfile /var/log/camilladsp.log
