# etc/systemd/system/camilladsp.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/camilladsp -p 1234 -o /var/log/camilladsp.log /srv/http/data/camilladsp/configs/camilladsp.yml
