# etc/systemd/system/camilladsp.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/camilladsp --port 1234 /srv/http/data/camilladsp/configs/camilladsp.yml
