# etc/systemd/system/bluezdbus.service
[Unit]
Description=Dbus for Bluetooth audio receiver
BindsTo=bluealsa-aplay.service

[Service]
ExecStart=/srv/http/bash/bluezdbus.py
