# etc/systemd/system/bluetoothbutton.service
[Unit]
Description=Bluetooth buttons

[Service]
Type=simple
ExecStart=/srv/http/bash/bluetoothbutton.sh
