# etc/systemd/system/bluealsa-aplay.service
[Unit]
Description=Bluetooth audio receiver

[Service]
ExecStart=bluealsa-aplay 00:00:00:00:00:00
