# etc/systemd/system/bluezdbus.service
[Unit]
Description=Dbus for Bluetooth audio receiver

[Service]
ExecStart=/srv/http/bash/bluezdbus.py
