# etc/systemd/system/camillagui.service
[Unit]
Description=CamillaDSP GUI

[Service]
Type=idle
ExecStart=/usr/bin/python /srv/http/settings/camillagui/main.py
