# etc/systemd/system/radio.service
[Unit]
Description=WebRadio metadata

[Service]
Type=simple
ExecStart=/srv/http/bash/status-radio.sh
