# etc/systemd/system/rotarymute.service
[Unit]
Description=Rotary mute
Requires=rotaryencoder.service
After=rotaryencoder.service

[Service]
ExecStart=/srv/http/bash/rotarymute.sh
