# etc/systemd/system/rotaryencoder.service
[Unit]
Description=Rotary encoder
BindsTo=rotarymute.service

[Service]
ExecStart=/srv/http/bash/rotaryencoder.sh
