# etc/systemd/system/bluealsa.service.d/override.conf
[Unit]
After=bluetooth.service
Requires=bluetooth.service
